# Deployment/Cloud/ansible/roles/private_egress/templates/private-egress.service.j2
[Unit]
Description={{ app_name }} private cloud egress route
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/{{ app_name }}-private-egress
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
